# services/idle-reporter.conf
[program:idle-reporter]
command=/usr/local/bin/idle-reporter.sh
autostart=true
autorestart=true
priority=50
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
